# packages/ipfs-kit-py/ipfs_kit_py-0.3.0-py3-none-any.whl/ipfs_kit_py/ipfs-cluster-follow.service
[Unit]
Description=IPFS Cluster Follow Daemon
After=network.target

[Service]
ExecStart=
Restart=on-failure
User=root
Group=root

[Install]
WantedBy=multi-user.target
