# packages/ipfs-kit-py/ipfs_kit_py-0.3.0.tar.gz/ipfs_kit_py-0.3.0/ipfs_kit_py/ipfs-cluster.service
[Unit]
Description=IPFS Cluster Daemon
After=network.target

[Service]
ExecStart=ipfs-cluster-service daemon
Restart=on-failure
User=root
Group=root

[Install]
WantedBy=multi-user.target
